# scripts/ctrlchecks-credential-service.service
[Unit]
Description=CtrlChecks Credential Service
Documentation=https://github.com/ctrlchecks/ctrlchecks-ai-workflow-os
After=network.target ctrlchecks-worker.service
Wants=ctrlchecks-worker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/ctrlchecks-credential-service
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/opt/ctrlchecks-credential-service/.env
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Phase 1 scaffold — minimal footprint (vault + OAuth in Phase 2)
MemoryMax=1G
LimitNOFILE=65536

# Stdout/stderr go to journald (sudo journalctl -u ctrlchecks-credential-service)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ctrlchecks-credential-service

[Install]
WantedBy=multi-user.target
