# scripts/ctrlchecks-trigger-service.service
[Unit]
Description=CtrlChecks Trigger Service
Documentation=https://github.com/ctrlchecks/ctrlchecks-ai-workflow-os
After=network.target ctrlchecks-worker.service
Wants=ctrlchecks-worker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/ctrlchecks-trigger-service
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/opt/ctrlchecks-trigger-service/.env
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Phase 1 scaffold — minimal footprint (real trigger handlers in Phase 2+)
MemoryMax=512M
LimitNOFILE=65536

# Stdout/stderr go to journald (sudo journalctl -u ctrlchecks-trigger-service)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ctrlchecks-trigger-service

[Install]
WantedBy=multi-user.target
